# etc/local.d/10-brightness.start
#!/bin/bash
MS=`cat /sys/class/backlight/intel_backlight/max_brightness`
LEVEL=0.6
VALUE=`echo "($MS*$LEVEL) / 1" | bc`

echo $VALUE > /sys/class/backlight/intel_backlight/brightness
